# rc.local
#!/bin/bash

bash /home/owen/boot_run_distcc.sh
